# crontab
# m h dom mon dow user	command
0 */1   * * * /home/sh/dotfiles/rclone.sh >/dev/null 2>&1 >> /etc/crontab
0 22    * * * docker restart keep_api_python_1 
